# guake-monitor-hotplug.service
[Unit]
Description=Clear Guake monitor cache on display hotplug

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'rm -f /run/user/*/guake-monitor/*.idx'
